# static/usr/lib/systemd/system/plymouth-reboot.service.d/core-override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/plymouthd --mode=reboot --attach-to-session --ignore-serial-consoles
