# static/usr/lib/systemd/system/secureboot-db.service.d/core-override.conf
[Service]
# We need the "exit 2" as the upstream service defines SuccessExitStatus=1
ExecCondition=sh -c 'test $(/usr/libexec/core/get-mode mode) = install || exit 2'
