# static/usr/lib/systemd/system/console-conf@.service.d/fix.conf
[Service]
SuccessExitStatus=22
ExecStopPost=
ExecStopPost=/usr/share/subiquity/console-conf-check-error-fix getty@%i.service
