# static/usr/lib/systemd/system/boot-efi.mount
[Unit]
RequiresMountsFor=/run/mnt/ubuntu-seed

[Mount]
Where=/boot/efi
What=/run/mnt/ubuntu-seed
Type=none
Options=bind
