# static/usr/lib/systemd/system/etc-shadow.mount
[Mount]
Where=/etc/shadow
What=/run/snapd/hybrid-users/shadow
Type=none
Options=bind
