# static/usr/lib/systemd/system/etc-group.mount
[Mount]
Where=/etc/group
What=/run/snapd/hybrid-users/group
Type=none
Options=bind
