# static/usr/lib/systemd/system/etc-gshadow.mount
[Mount]
Where=/etc/gshadow
What=/run/snapd/hybrid-users/gshadow
Type=none
Options=bind
