# static/usr/lib/systemd/system/etc-passwd.mount
[Mount]
Where=/etc/passwd
What=/run/snapd/hybrid-users/passwd
Type=none
Options=bind
