# static/usr/lib/systemd/system/fix-hctosys.service
[Unit]
Description=Work around automatic hctosys

[Service]
Type=oneshot
ExecStart=/usr/lib/core/fix-hctosys
