# static/usr/lib/systemd/system/plymouth-show-install.service
[Unit]
Description=Display state on splash
After=plymouth-start.service
ConditionPathExists=/run/plymouth/pid

[Service]
ExecCondition=sh -c 'test $(/usr/libexec/core/get-mode mode) = install'
ExecStart=-/bin/plymouth display-message --text='Installing Ubuntu Core...'
Type=oneshot
TimeoutSec=20
