# static/usr/lib/systemd/system/system-systemd\x2dfsck.slice
# We need to mark this slice as not conflicting with shutdown.target.
# Otherwise it would drag mounts to be unmounted.
[Unit]
Description=Fsck Units Slice
DefaultDependencies=no

