# static/usr/lib/systemd/system/var-lib-optee\x2dclient-data-tee.mount
[Unit]
ConditionPathIsMountPoint=/run/mnt/tee-data
DefaultDependencies=no

[Mount]
What=/run/mnt/tee-data
Where=/var/lib/optee-client/data/tee
Options=bind
Type=none
