# static/usr/lib/systemd/system/plymouth-quit.service.d/core-override.conf
[Unit]
After=snapd.seeded.service
ConditionKernelCommandLine=snapd_recovery_mode=run

[Service]
ExecStartPost=-systemctl stop splash-client.service
