# static/usr/lib/systemd/system/plymouth-quit.service.d/core-override.conf
[Unit]
After=snapd.seeded.service

[Service]
ExecCondition=sh -c 'test $(/usr/libexec/core/get-mode mode) = run'
ExecStartPost=-systemctl stop splash-client.service
