# static/usr/lib/systemd/system/systemd-hostnamed.service.d/writable.conf
[Service]
Environment=SYSTEMD_ETC_HOSTNAME=/etc/writable/hostname
Environment=SYSTEMD_ETC_MACHINE_INFO=/etc/writable/machine-info
