# static/usr/lib/systemd/system/systemd-timedated.service.d/writable.conf
[Service]
Environment=SYSTEMD_ETC_LOCALTIME=/etc/writable/localtime
Environment=SYSTEMD_ETC_ADJTIME=/etc/writable/adjtime
