# static/usr/lib/systemd/system/plymouth-start.service.d/core-override.conf
[Unit]
Wants=plymouth-quit.service splash-client.service
Wants=plymouth-show-run.service plymouth-show-install.service
Wants=plymouth-show-recover.service plymouth-show-factory-reset.service

[Service]
ExecStart=
ExecStart=/usr/sbin/plymouthd --mode=boot --pid-file=/run/plymouth/pid --attach-to-session --ignore-serial-consoles
