# static/usr/lib/systemd/system/plymouth-poweroff.service.d/core-override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/plymouthd --mode=shutdown --attach-to-session --ignore-serial-consoles
