# static/usr/lib/systemd/system/tee-supplicant.service.d/core-override.conf
# SPDX-License-Identifier: BSD-2-Clause
# Adapted from https://github.com/OP-TEE/optee_client/blob/6486773583b5983af8250a47cf07eca938e0e422/tee-supplicant/tee-supplicant%40.service.in
[Unit]
DefaultDependencies=no
BindsTo=dev-teepriv0.device var-lib-optee\x2dclient-data-tee.mount
After=dev-teepriv0.device var-lib-optee\x2dclient-data-tee.mount

[Service]
Type=notify
ExecStartPre=-/usr/sbin/modprobe -v -r tpm_ftpm_tee
ExecStartPost=-/usr/sbin/modprobe -v tpm_ftpm_tee
ExecStop=-/bin/sh -c "/sbin/modprobe -v -r tpm_ftpm_tee ; /bin/kill $MAINPID"
